# files/etc/hotplug.d/usb/20-fm350-usbids
#!/bin/sh
# FM350-GL USB ID Registration Hotplug Script
# POSIX/ash compliant
# v22.1: Relaxed filter to e8d/712[567]/.*

# Rate limit >= 60s
stamp="/tmp/fm350.usbids.stamp"
now=$(date +%s)
if [ -f "$stamp" ]; then
  last=$(cat "$stamp" 2>/dev/null)
  [ -n "$last" ] && [ $((now - last)) -lt 60 ] && exit 0
fi
echo "$now" > "$stamp"

# Filter: only add/bind actions
case "$ACTION" in
  add|bind)
    ;;
  *)
    exit 0
    ;;
esac

# Match FM350 VID:PID (0e8d:7125/7126/7127, any interface)
# v22.1: Changed from /1 to /.*
echo "$PRODUCT" | grep -qE '^e8d/712[567]/.*' || exit 0

# Register USB IDs to option driver (fallback to generic)
echo "0e8d 7127 ff" > /sys/bus/usb-serial/drivers/option1/new_id 2>/dev/null || \
echo "0e8d 7127"    > /sys/bus/usb-serial/drivers/generic/new_id  2>/dev/null

echo "0e8d 7126 ff" > /sys/bus/usb-serial/drivers/option1/new_id 2>/dev/null || \
echo "0e8d 7126"    > /sys/bus/usb-serial/drivers/generic/new_id  2>/dev/null

echo "0e8d 7125 ff" > /sys/bus/usb-serial/drivers/option1/new_id 2>/dev/null || \
echo "0e8d 7125"    > /sys/bus/usb-serial/drivers/generic/new_id  2>/dev/null

logger -t fm350 -p daemon.info "USB IDs registered for FM350-GL (action=$ACTION, product=$PRODUCT)"

exit 0
