# files/etc/hotplug.d/usb/60-fm350
#!/bin/sh
# FM350-GL Manager Watchdog Hotplug Script
# POSIX/ash compliant

# Rate limit >= 60s
stamp="/tmp/fm350.watchdog.stamp"
now=$(date +%s)
if [ -f "$stamp" ]; then
  last=$(cat "$stamp" 2>/dev/null)
  [ -n "$last" ] && [ $((now - last)) -lt 60 ] && exit 0
fi
echo "$now" > "$stamp"

# Filter: only add/bind actions for FM350
case "$ACTION" in
  add|bind)
    ;;
  *)
    exit 0
    ;;
esac

echo "$PRODUCT" | grep -qE '^e8d/712[567]/.*' || exit 0

# Check if manager is running
if ! pidof fm350-manager >/dev/null 2>&1; then
  logger -t fm350 -p daemon.warn "Manager not running, restarting service (one-time)"
  /etc/init.d/fm350-manager restart
fi

exit 0
