# files/etc/hotplug.d/usb/60-fm350
#!/bin/sh
# FM350-GL Manager Watchdog Hotplug Script
# POSIX/ash compliant

. /usr/lib/fm350/functions.sh

# Filter: only bind/add actions for FM350
[ "$ACTION" = "add" ] || [ "$ACTION" = "bind" ] || exit 0

# Match FM350 VID:PID
echo "$PRODUCT" | grep -qE '^e8d/(7127|7126|7125)/' || exit 0

# Rate limiting
STAMP_FILE="/tmp/fm350-hotplug-watchdog.stamp"
mkdir -p /tmp

if [ -f "$STAMP_FILE" ]; then
  LAST=$(cat "$STAMP_FILE" 2>/dev/null)
  NOW=$(date +%s)
  
  if [ -n "$LAST" ] && [ $((NOW - LAST)) -lt 60 ]; then
    exit 0
  fi
fi

date +%s > "$STAMP_FILE"

# Check if manager is running
if ! pidof fm350-manager >/dev/null 2>&1; then
  log_limited warn "Hotplug: Manager not running, restarting service"
  /etc/init.d/fm350-manager restart
fi

exit 0